# sysutils/pkg-updated/pkg-plist
bin/pkg-updated
etc/pkg-updated/pkg-updated.conf.sample
etc/rc.d/pkg-updated
@dir /var/db/pkg-updated/archive
@dir /var/db/pkg-updated
